# services/prediction.service
[Unit]
Description=run stem club prediction
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/stem_club/
ExecStart=/home/pi/stem_club/predict.sh --model=0 --source=1
RemainAfterExit=yes
User=pi
GuessMainPID=no

[Install]
WantedBy=multi-user.target
